# debian/ailab-web.service
[Unit]
Description=AI Lab web management interface
Documentation=https://github.com/lemonade-sdk/ailab
After=network.target

[Service]
ExecStart=/usr/bin/ailab web
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
